# conf/rc.local
#!/usr/bin/env bash

xmodmap /home/jakob/.xmodmap
